# raspberry/config_files/etc_supervisor_conf.d_aquarium.conf
[program:aquarium]
command=/home/corneliu/env/bin/python /home/corneliu/env/aquarium/raspberry/server.py
logfile=/home/corneliu/env/log.log
logfile_maxbytes=50MB
